# mirror/fhlun-core-mirror.timer
[Unit]
Description=Daily FHLUN core mirror synchronization

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
