# mirror/fhlun-core-mirror.service
[Unit]
Description=Synchronize FHLUN core mirror assets
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/fhlun-core-mirror-sync
